# examples/production-host/systemd/tinyidp.service
[Unit]
Description=tinyidp strict production host
Documentation=https://github.com/go-go-golems/tiny-idp
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=tinyidp
Group=tinyidp
StateDirectory=tinyidp
LogsDirectory=tinyidp
UMask=0077
ExecStart=/usr/local/bin/tinyidp serve-production --addr :8443 --listener-mode direct-tls --issuer https://idp.example.com --clients-file /etc/tinyidp/catalog/clients.json --theme-dir /etc/tinyidp/themes --theme-catalog-file /etc/tinyidp/themes/themes.json --signup-program-file /etc/tinyidp/signup.js --db /var/lib/tinyidp/idp.db --audit-path /var/log/tinyidp/audit.jsonl --token-secret-file /etc/tinyidp/token-secret --tls-cert /etc/tinyidp/tls.crt --tls-key /etc/tinyidp/tls.key
Restart=on-failure
RestartSec=5s
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
ReadWritePaths=/var/lib/tinyidp /var/log/tinyidp
ReadOnlyPaths=/etc/tinyidp
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
AmbientCapabilities=CAP_NET_BIND_SERVICE
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
